# deploy/systemd/nous-auto-update.service
[Unit]
Description=Update nous to origin/main and restart the service
Documentation=https://github.com/rmednitzer/nous
After=network-online.target nous.service
Wants=network-online.target
# We do not Requires=nous.service: this unit can run before nous.service
# is healthy (in fact it may be what fixes a broken nous.service).

[Service]
Type=oneshot
ExecStart=/opt/nous/deploy/auto-update.sh
# install.sh can take a minute on a fresh dependency set; cap generously.
TimeoutStartSec=10m
StandardOutput=journal
StandardError=journal
# AUDIT-2026-05-20 H8: explicit no-restart + clean exit status.
# Restart=no is the default for Type=oneshot but stating it inline
# documents the intent: a failed update does NOT loop within the
# unit; the timer fires the next tick, and auto-update.sh checks
# /var/log/nous/auto-update.last_failed to skip a commit that
# already broke. The pair (Restart=no + last_failed skip) replaces
# the implicit every-five-minutes retry behaviour with explicit
# refusal until an operator intervenes via
# deploy/auto-update-rollback.sh.
Restart=no
SuccessExitStatus=0
# Auto-update touches /opt/nous, /etc/nous, /var/log/nous, /etc/systemd,
# and restarts the unit; runs as root by necessity.
