# deploy/systemd/nous-state-flush.service
[Unit]
Description=nous daily state flush
After=nous.service

[Service]
Type=oneshot
User=nous
Group=nous
EnvironmentFile=/etc/nous/nous.env
WorkingDirectory=/opt/nous
# `flush` truncates the SQLite WAL (single PRAGMA wal_checkpoint(TRUNCATE))
# and fsyncs the audit handler. A single daily checkpoint bounds disk
# thrashing while keeping audit fidelity tight (every tool call is
# already fsynced individually -- the daily run is the WAL safety net).
ExecStart=/opt/nous/venv/bin/python -m nous flush
TimeoutStartSec=30
