# deploy/systemd/nous.service
[Unit]
Description=nous simulator MCP server
Documentation=https://github.com/rmednitzer/nous
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nous
Group=nous
EnvironmentFile=/etc/nous/nous.env
WorkingDirectory=/opt/nous
ExecStart=/opt/nous/venv/bin/python -m nous serve
# Pre-stop flush keeps the audit log durable across SIGTERM: the
# AuditLogger fsyncs after every record, but the SQLite WAL needs an
# explicit checkpoint before the process exits to keep the on-disk
# snapshot consistent.
ExecStop=/opt/nous/venv/bin/python -m nous flush
TimeoutStopSec=20
Restart=on-failure
RestartSec=5s
# The audit log lives at /var/log/nous/audit.jsonl; logrotate manages it.
ReadWritePaths=/var/lib/nous /var/log/nous

# Hardening (24.04-baseline directives).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
LockPersonality=true

# Hardening added for the Ubuntu 26.04 / systemd 256+ baseline (ADR 0016).
# Older systemd versions ignore unknown directives with a warning; the
# service still starts on 24.04.
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
MemoryDenyWriteExecute=true
RemoveIPC=true
KeyringMode=private
UMask=0077
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @debug @mount @cpu-emulation @obsolete @raw-io @reboot @swap
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
